# power_monitor.service
[Unit]
Description=Power Monitor Push Server
After=network.target

[Service]
User=root
WorkingDirectory=/root/geminicli/light-monitor-kyiv
EnvironmentFile=/root/geminicli/light-monitor-kyiv/.env
ExecStart=/root/geminicli/light-monitor-kyiv/venv/bin/python power_monitor_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
